# config/supervisord.conf
[supervisord]
nodaemon=true

[program:payara]
command=/opt/payara41/glassfish/bin/asadmin start-domain -v 

